# src/template.socket
[Unit]
Description=FAUST CTF helloworld service

[Socket]
ListenStream=666
Accept=true

[Install]
WantedBy=sockets.target
